# src/default/containers.manifest
# This file is parsed by containers.initd => /etc/init.d/containers

selfhostedpro/yacht:latest /bin/sleep 60
